# config/includes.chroot/etc/xdg/autostart/volti.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.2.3-7
Type=Application
Name=volti
Comment=
Exec=volti
StartupNotify=false
Terminal=false
Hidden=false
